# business/supervisor.conf
; ============================================================
; 业务应用 - Supervisor 配置（business/app.py）
; ============================================================
; 【必改】directory、command 改为实际路径
; ============================================================

[program:business]
command=/path/to/python_api/venv/bin/python portal/run_business.py
directory=/path/to/python_api

user=root
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=10

environment=PYTHONUNBUFFERED="1",BUSINESS_PORT="5003"

redirect_stderr=true
stdout_logfile=/var/log/supervisor/business.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
