# portal/supervisor_data_map.conf
; ============================================================
; 数据地图模块 - 单端口 5001（数据字典 + 数据血缘）
; ============================================================
;
; 【必改】directory、command、stdout_logfile 改为实际路径
; 使用：supervisorctl reread && supervisorctl update && supervisorctl start data_map
; ============================================================

[program:data_map]
directory=/path/to/python_api
command=/path/to/python_api/venv/bin/python portal/run_data_map.py

user=root
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=10

environment=PYTHONUNBUFFERED="1",DATA_MAP_PORT="5001"

redirect_stderr=true
stdout_logfile=/var/log/supervisor/data_map.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
