# portal/supervisor_run_all.conf
; ============================================================
; 单进程单端口（run_all.py）- 门户 + 数据地图 + 技术支持 + 业务 全在一个进程
; ============================================================
;
; 与「四进程」二选一：用本配置则只起 portal，不要起 data_map、support、business。
; 【必改】directory、command、stdout_logfile 改为实际路径。
; ============================================================

[program:portal_run_all]
directory=/path/to/python_api
command=/path/to/python_api/venv/bin/python portal/run_all.py

user=root
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=10

environment=PYTHONUNBUFFERED="1"

redirect_stderr=true
stdout_logfile=/var/log/supervisor/portal_run_all.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
