# portal/supervisor_support.conf
; ============================================================
; 技术支持模块 - 单端口 5002（Excel 入库 Hive，可扩展）
; ============================================================
;
; 【必改】directory、command、stdout_logfile 改为实际路径
; 使用：supervisorctl reread && supervisorctl update && supervisorctl start support
; ============================================================

[program:support]
directory=/path/to/python_api
command=/path/to/python_api/venv/bin/python portal/run_support.py

user=root
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=10

environment=PYTHONUNBUFFERED="1",SUPPORT_PORT="5002"

redirect_stderr=true
stdout_logfile=/var/log/supervisor/support.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
