# portal/supervisor_portal.conf
; ============================================================
; 统一门户（仅门户，四进程部署时用）
; ============================================================
;
; 四进程部署：portal(5000) + data_map(5001) + support(5002) + business(5003) 互不影响。
; 【必改】directory、command、stdout_logfile 改为实际路径。
; 四进程时门户 iframe 指向三个独立端口，需设置 DATA_MAP_URL、SUPPORT_URL、BUSINESS_URL、PORTAL_URL
; （见下方 environment 示例，把 192.168.168.219 改成你的服务器 IP 或域名）。
;
; 使用：supervisorctl reread && supervisorctl update && supervisorctl start portal
; ============================================================

[program:portal]
directory=/path/to/python_api
command=/path/to/python_api/venv/bin/python -m portal.app

user=root
autostart=true
autorestart=true
startretries=3
startsecs=3
stopwaitsecs=10

; 四进程时门户需知道三个模块的地址，以及门户自身地址（子模块登录/登出重定向用）
environment=PYTHONUNBUFFERED="1",PORTAL_PORT="5000",DATA_MAP_URL="http://192.168.168.219:5001",SUPPORT_URL="http://192.168.168.219:5002",BUSINESS_URL="http://192.168.168.219:5003",PORTAL_URL="http://192.168.168.219:5000"

redirect_stderr=true
stdout_logfile=/var/log/supervisor/portal.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
